# board/raspberrypi/fs_edgemap_display/etc/udev/rules.d/90-nk-macsec.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="20a0", ATTR{idProduct}=="4109", RUN+="/opt/nk-macsec/rekey.sh"

